# hypr/.hypr/s6-usertree/sv/dbus/run
#!/bin/bash
exec dbus-daemon --nofork --session --address=$DBUS_SESSION_BUS_ADDRESS 
